# nginx/cron.d/pass
@reboot root mkdir -p /dev/shm/pass/share && chown -R root:nginx /dev/shm/pass && chmod -R 750 /dev/shm/pass
0 * * * * root find /dev/shm/pass/share/ -type f -mmin +1439 -delete > /dev/null 2>&1
* * * * * root find /dev/shm/pass/share/ -exec stat -c \%n#\%X#\%Y {} \; | awk -F# '{if ($2 != $3) print $1}' | xargs /bin/rm -f 2>/dev/null
